# systemd/OPTIMAL_BOT_TG.service
[Unit]
Description=gpt-bot
After=multi-user.target

[Service]
User=student
Group=student
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/student/OPTIMAL_BOT_TG/main.py

[Install]
WantedBy=multi-user.target